# nftables/nftables.service.d/override.conf
[Service]
ExecStart=
ExecStart=/sbin/nft -f /etc/nftables/nftables.conf
