# tinc/tinc.service.d/override.conf
[Service]
ExecStartPre=/usr/bin/bash -c 'if [ ! -e /dev/net/tun ]; then mkdir -p /dev/net; mknod /dev/net/tun c 10 200; fi'
